# firmware/agent/systemd/clawos-agent-bin.service
[Unit]
Description=CLAWOS Raspberry Pi Agent (Binary)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/clawos/agent
ExecStart=/opt/clawos/agent/clawos-agent-linux-arm64
Restart=always
RestartSec=3
User=clawos
Group=clawos
Environment=AGENT_CONFIG_PATH=/opt/clawos/agent/config/agent.config.json

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/clawos/agent/runtime /opt/clawos/agent/config

[Install]
WantedBy=multi-user.target
